# sense-pulse.service
[Unit]
Description=Sense Pulse - Pi-hole and Tailscale Status Display with Web Dashboard
After=network.target pihole-FTL.service tailscaled.service
Wants=network.target

[Service]
Type=simple
ExecStart=/home/pi/.local/bin/uv run --project /home/pi/Personal/sense-pulse sense-pulse
WorkingDirectory=/home/pi/Personal/sense-pulse
Restart=on-failure
RestartSec=10

User=pi
Group=pi

Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/pi/.local/bin"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sense-pulse

[Install]
WantedBy=multi-user.target
